# config/crontab
0   0   *   *   *   /var/lib/dubsar/current/script/wotd
